# systemd/blofin-stack-ingestor.service
[Unit]
Description=Blofin monitoring ingestor (websocket -> sqlite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/rob/.openclaw/workspace/blofin-stack
EnvironmentFile=/home/rob/.openclaw/workspace/blofin-stack/.env
ExecStart=/home/rob/.openclaw/workspace/blofin-stack/.venv/bin/python /home/rob/.openclaw/workspace/blofin-stack/ingestor.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
